# resistor-station.service
[Unit]
Description=Resistor Station Teaching App
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=azmeares
WorkingDirectory=/home/azmeares/Resistor-Station/pi-app
ExecStart=/home/azmeares/Resistor-Station/pi-app/start.sh
Restart=on-failure
RestartSec=3

# Give hardware time to initialise on boot
ExecStartPre=/bin/sleep 2

# Allow access to I2C, serial
SupplementaryGroups=video input i2c dialout gpio

[Install]
WantedBy=multi-user.target
