# resistor-station.service
[Unit]
Description=Resistor Station Teaching App
After=display-manager.service
Wants=display-manager.service

[Service]
Type=simple
User=azmeares
WorkingDirectory=/home/azmeares/Resistor-Station/pi-app
ExecStart=/home/azmeares/Resistor-Station/pi-app/start.sh
Restart=on-failure
RestartSec=3

# Give hardware time to initialise on boot
# Wait for display manager + hardware to be ready
ExecStartPre=/bin/sleep 5

# Allow access to I2C, serial
SupplementaryGroups=video input i2c dialout gpio

[Install]
WantedBy=graphical.target
